# jpackage/resources/postinst
#!/bin/sh

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
    	useradd --system --quiet --home /opt/APPLICATION_PACKAGE \
			--comment "User account for running the APPLICATION_PACKAGE" fop-pdf-user

		echo APPLICATION_DIRECTORY
		cp /opt/APPLICATION_PACKAGE/install -m 644 %{name}.service %{buildroot}%{_unitdir}/%{name}.service			

		deb-systemd-helper enable 'APPLICATION_PACKAGE.service'
		deb-systemd-helper start 'APPLICATION_PACKAGE.service'
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
