# jpackage/resources/postinst
#!/bin/sh

systemctl enable APPLICATION_PACKAGE
systemctl start APPLICATION_PACKAGE

exit 0