# jpackage/resources/prerm
#!/bin/sh

systemctl stop APPLICATION_PACKAGE
systemctl disable APPLICATION_PACKAGE

exit 0